# firefox-kiosk.service
[Unit]
Description=Firefox Kiosk Mode
# Start after the graphical session is ready for the user
After=graphical-session.target
PartOf=graphical-session.target

[Service]
# User to run Firefox as (replace 'your_kiosk_user' with the actual username)
User=sysadmin
# Group for the user (optional, often the same as the username)
# Group=your_kiosk_user

# Environment variables needed for GUI applications
Environment=DISPLAY=:0
# For Wayland, you might need other variables like WAYLAND_DISPLAY
# Environment=XAUTHORITY=/home/your_kiosk_user/.Xauthority # May be needed depending on Xorg setup

# Command to start Firefox in kiosk mode
# Adjust the path to firefox if it's different on your system
# --no-remote is often useful to ensure a new instance is started
ExecStart=/usr/bin/firefox --kiosk http://localhost --no-remote

# Restart behavior
Restart=always
RestartSec=5

# Standard output and error logging
StandardOutput=journal
StandardError=journal

[Install]
# Start when the graphical session target is reached
WantedBy=graphical-session.target